# resources/chef-repo/site-cookbooks/rails_part/templates/default/puma_init_script.erb
#!/bin/bash
#
# <%= node['rails_part']['app']['name'] %>

# chkconfig: 2345 82 55
:# processname <%= node['rails_part']['app']['name'] %>
# description: Runs <%= node['rails_part']['app']['name'] %> for nginx integration.

# Include RedHat function library
. /etc/rc.d/init.d/functions

# The name of the service
NAME=<%= node['rails_part']['app']['name'] %>

# The username and path to the myapp source
USER=root
APP_PATH=<%= node['rails_part']['app']['path'] %>/current

# The PID and LOCK files used by puma and sidekiq
UPID=<%= node['rails_part']['app']['path'] %>/shared/puma/<%= node['rails_part']['app']['name'] %>.pid
ULOCK=/var/lock/subsys/$NAME

# The options to use when running puma
OPTS="-C <%= node['rails_part']['app']['path'] %>/shared/puma/<%= node['rails_part']['app']['name'] %>.config -e <%= node['rails_part']['app']['rails_env'] %>"

# Ruby related path update
RUBY_PATH_PATCH="PATH=$PATH:/opt/rbenv/shims:/usr/local/lib && export PATH && "
BUNDLE_CMD=bundle
PUMA_CMD=puma

start() {
  cd $APP_PATH

  # Start puma
  echo -n $"Starting $NAME: "
  daemon --pidfile=$UPID --user=$USER "$RUBY_PATH_PATCH $BUNDLE_CMD exec $PUMA_CMD $OPTS"
  puma=$?
  [ $puma -eq 0 ] && touch $ULOCK
  echo

  return $puma
}

stop() {
  cd $APP_PATH

  # Stop puma
  echo -n $"Stopping $NAME: "
  killproc -p $UPID

  puma=$?

  [ $puma -eq 0 ] && rm -f $ULOCK
  echo

  rm <%= node['rails_part']['app']['path'] %>/shared/puma/<%= node['rails_part']['app']['name'] %>.pid
  rm <%= node['rails_part']['app']['path'] %>/shared/puma/<%= node['rails_part']['app']['name'] %>.state
  rm <%= node['rails_part']['app']['path'] %>/shared/puma/<%= node['rails_part']['app']['name'] %>_control.sock


  return $puma
}

restart() {
  stop
  start
}

get_status() {
  status -p $UPID $NAME
}

query_status() {
  get_status >/dev/null 2>&1
}

case "$1" in
  start)
    query_status && exit 0
    start
    ;;
  stop)
    query_status || exit 0
    stop
    ;;
  restart)
    restart
    ;;
  status)
    get_status
    exit $?
    ;;
  *)
    N=/etc/init.d/$NAME
    echo "Usage: $N {start|stop|restart|status}" >&2
    exit 1
    ;;
esac

exit 0
